# external/rootfs_overlay/etc/cron/crontabs/root
0 3 * * * /usr/sbin/hwclock --systohc -f /dev/rtc0